# server/monfetch-server.service
[Unit]
Description=monfetch fetch-style web monitoring service
After=network.target

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/sbin/uwsgi monfetch-uwsgi.ini
WorkingDirectory=/usr/local/bin/monfetch/server

[Install]
WantedBy=multi-user.target
